# wignerf.service
[Unit]
Description=Wigner Function Simulator
#After=network-online.target nginx.service
#Requires=nginx.service

[Service]
Type=simple
User=tigran
Group=tigran

WorkingDirectory=/home/tigran/python/wignerf/
Environment="PYTHONUNBUFFERED=1"
# cupy JIT-compiles CUDA kernels at runtime (NVRTC) and caches the cubins.
# ProtectHome=read-only makes the default ~/.cupy unwritable, AND PrivateTmp
# wipes /tmp every restart — so the cache MUST be a PERSISTENT writable dir,
# else every start recompiles (~7 s/variant, measured; x4 variants serialized
# by the GIL = a minute+), blocking the first frame each time. CacheDirectory
# gives /var/cache/wignerf (systemd-owned, RW in the sandbox, survives
# restarts): compile once, sub-second thereafter. No-op on the CPU-only VPS.
CacheDirectory=wignerf
Environment=CUPY_CACHE_DIR=/var/cache/wignerf/cupy
# Same story for matplotlib (imported at startup by the mp4 exporter): the
# read-only home makes ~/.config/matplotlib unwritable, so it falls back to a
# fresh /tmp dir that PrivateTmp then wipes — rebuilding the font cache on
# EVERY start and logging a warning about it. Measured here: 1.12 s cold vs
# 0.39 s warm for import + first draw, paid again on every --reload in dev.
# Lives beside CUPY_CACHE_DIR rather than in wignerf.env because the path only
# exists by virtue of CacheDirectory= three lines up; splitting them lets one
# be changed without the other. matplotlib creates the subdir if missing.
Environment=MPLCONFIGDIR=/var/cache/wignerf/matplotlib
EnvironmentFile=/home/tigran/python/wignerf/wignerf.env
ExecStart=/home/tigran/python/wignerf/start.sh

Restart=always
RestartSec=5

# --- Security Constraints ---
NoNewPrivileges=yes
ProtectSystem=strict
#ReadWritePaths=
ProtectHome=read-only
PrivateTmp=yes
# PrivateDevices MUST stay off: =yes replaces /dev with a minimal private
# tmpfs that has NO /dev/nvidia* nodes, so cupy sees zero CUDA devices and the
# solver silently falls back to CPU. DeviceAllow cannot re-add nodes under
# PrivateDevices=yes, so leave it off to expose the GPUs. (No-op on the VPS.)
PrivateDevices=no
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
